# meta-citadel/recipes-citadel/citadel-config/files/systemd/launch-default-realm.path
[Unit]
Description=Wait for creation of /run/user/1000/pulse/native

[Path]
PathChanged=/run/user/1000/pulse/native
