# meta-citadel/recipes-citadel/bbappends/networkmanager/watch-resolvconf.path
[Unit]
Description=Watch /run/NetworkManager

[Path]
PathChanged=/run/NetworkManager
MakeDirectory=true

[Install]
WantedBy=multi-user.target
